# src/Scripts/dev/fabien_compiegne.cron
# $Id$
# check on commit runs on every minute
* * * * * /home/fabien/CHECK/pips_check_compile --source /home/fabien/CHECK/pipsrc.sh pips-commit@cri.ensmp.fr /home/fabien/CHECK/prod /home/fabien/CHECK/validation
# force once a day
7 0 * * * /home/fabien/CHECK/pips_check_compile --wait --force --report --source /home/fabien/CHECK/pipsrc.sh pips-commit@cri.ensmp.fr /home/fabien/CHECK/prod /home/fabien/CHECK/validation
# test autoconf version
@hourly /home/fabien/AUTOC/pips_check_compile --autoconf --source=/home/fabien/AUTOC/pipsrc.sh pips-commit@cri.ensmp.fr /home/fabien/AUTOC/prod /home/fabien/AUTOC/validation
# force once a day
7 1 * * * /home/fabien/AUTOC/pips_check_compile --wait --force --report --autoconf --source=/home/fabien/AUTOC/pipsrc.sh pips-commit@cri.ensmp.fr /home/fabien/AUTOC/prod /home/fabien/AUTOC/validation
# 32 bit cross compilation
7 2 * * * /home/fabien/CROSS/pips_check_compile --force --report --source=/home/fabien/CROSS/pipsrc.sh coelho@cri.ensmp.fr /home/fabien/CROSS/prod
