# src/Scripts/dev/fabien_riquet.cron
# $Id$
#
# Debian Lenny i686
#
# riquet is already svn & trac & redmine
# so just force compilation & validation once a day
25 2 * * * /home/fabien/INT32/pips_check_compile --force --report --source /home/fabien/INT32/pipsrc.sh coelho@cri.ensmp.fr /home/fabien/INT32/prod /home/fabien/INT32/validation
