# src/Scripts/dev/fabien_compiegne.cron
# $Id$
# use manually install fortran & runtime
# for debian/ubuntu compatibility in some tests
PATH=/usr/local/bin:/usr/bin:/bin
LD_LIBRARY_PATH=/usr/local/lib64
#PIPS_VALIDATION_SKIP = Pyps
PIPS_VALIDATION_NO_PYPS=1
PIPS_F77=gfortran
PIPS_F90=gfortran
* * * * * /home/fabien/pips_check_compile pips-commit@cri.ensmp.fr /home/fabien/CHECK/prod /home/fabien/CHECK/validation
@daily /home/fabien/pips_check_compile --wait --force pips-commit@cri.ensmp.fr /home/fabien/CHECK/prod /home/fabien/CHECK/validation
