# src/Scripts/dev/pips_hendaye.cron
# $Id$
# temporary fix for pips@verone because svn version 1.5 is required
30 2 * * * /usr/bin/svn up $HOME/PIPS/prod $HOME/PIPS/validation > $HOME/PIPS/svn.out 2> $HOME/PIPS/svn.err
