# packages/PIPS/pips/src/Scripts/dev/fabien_compiegne.cron
# $Id$
# check on commit runs on every minute
* * * * * /home/fabien/CHECK/pips_check_compile --source /home/fabien/CHECK/pipsrc.sh pips-commit@cri.ensmp.fr /home/fabien/CHECK/prod /home/fabien/CHECK/validation
# force once a day
@daily /home/fabien/CHECK/pips_check_compile --wait --force --source /home/fabien/CHECK/pipsrc.sh pips-commit@cri.ensmp.fr /home/fabien/CHECK/prod /home/fabien/CHECK/validation
# test autoconf version
@daily /home/fabien/AUTO/pips_check_compile --force --report --autoconf --source=/home/fabien/AUTO/pipsrc.sh coelho@cri.ensmp.fr /home/fabien/AUTO/prod /home/fabien/AUTO/validation
